# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=true

[program:beanstalkd]
command=beanstalkd -l 127.0.0.1 -p 11301

[program:aurora]
command=aurora -c /home/aurora.toml
