# supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
nodaemon=true

[program:beanstalkd]
command=beanstalkd

[program:aurora]
command=aurora -c /home/aurora.toml
